# deploy/4sq.conf
[program:4sq]
directory = /home/ubuntu/4squaredaway/
command = /home/ubuntu/4squaredaway/venv/bin/gunicorn -c deploy/gunicorn_config.py fsqaway.wsgi:app
user = ubuntu
